# lib/etc/arnold_gt.conf
# ArnoldGT - Process Shelby GT jobs, looking for video and adding to the DBs

description "ArnoldGT - The 'link processor' for ShelbyGT"
author "Dan Spinosa <dan@shelby.tv>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# Start when networking has been started
start on started network-services

# And stop when networking is stopping
stop on stopping network-services

# Automatically restart process if crashed (not limiting how many times to restart)
#respawn

# Essentially lets upstart know the process will detach itself to the background
expect fork

# We aren't actually running Arnold directly, we're running a supervisor process that keeps X arnolds running (to take advantage of our cores)
# `sudo -u gt` changes to our user and `-i` sets up the environment (needed for rvm)
exec sudo -u gt -i /home/gt/shelby_gt/current/lib/gt/arnold/pump_iron.rb 7 /home/gt/shelby_gt/current/lib/gt/arnold/arnold_gt.rb --machine_name=gtX  > /home/gt/shelby_gt/shared/log/pump_iron.log 2>&1 &
