# lib/etc/mongo_mms.conf
# Ubuntu upstart file at /etc/init/mongo_mms.conf
# Starts the MMS agent on gt-api-a

pre-start script
  mkdir -p /home/gt/mms-agent/log
end script

start on runlevel [345]
stop on runlevel [06]

script
  # running as GT since we use that users home directory
  exec nohup /usr/bin/env python /home/gt/mms-agent/agent.py >> /home/gt/mms-agent/log/agent.log 2>&1;
end script