# lib/etc/mongo_all_arbs.conf
# Ubuntu upstart file at /etc/init/mongo_all_arbs.conf
# This script runs on gt-api-a
# It starts all the arbiter instances

limit nofile 20000 20000

pre-start script
    mkdir -p /var/log/mongodb/
    mkdir -p /var/lib/mongodb/
    mkdir -p /var/lib/mongodb/gtDashboard
    mkdir -p /var/lib/mongodb/gtVideo
    mkdir -p /var/lib/mongodb/gtConversation
    mkdir -p /var/lib/mongodb/gtRollFrame
    mkdir -p /var/lib/mongodb/gtUserAction
    mkdir -p /var/lib/mongodb/gtDeeplinkCache
    #mkdir -p /var/lib/mongodb/shelbySet
end script

start on runlevel [2345]

task

script
  start mongo_arb PORT=27019 REPLSET="gtDashboard"
  start mongo_arb PORT=27020 REPLSET="gtVideo"
  start mongo_arb PORT=27021 REPLSET="gtConversation"
  start mongo_arb PORT=27022 REPLSET="gtRollFrame"
  start mongo_arb PORT=27023 REPLSET="gtUserAction"
  start mongo_arb PORT=27024 REPLSET="gtDeeplinkCache"
  #start mongo_arb PORT=27021 REPLSET="shelbySet"
end script