# lib/etc/nginx.conf
# Ubuntu upstart file at /etc/init/nginx.conf

start on runlevel [2345]
stop on runlevel [06]

script
  exec start-stop-daemon --start --quiet --exec  /opt/nginx/sbin/nginx;
end script